# elonbot2022.conf
[program:elonbot22]
directory=/home/ubuntu/bots/elonbot/elonbot-2022
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/elonbot/elonbot-2022/logfile_err.log
stdout_logfile=/home/ubuntu/bots/elonbot/elonbot-2022/logfile.log